# templates/supervisord.ini
[unix_http_server]
file=${parts.buildout.directory}/var/supervisor.sock                       ; path to your socket file

[supervisord]
logfile=${parts.buildout.directory}/var/log/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=${parts.buildout.directory}/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=${parts.buildout.directory}/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://${parts.buildout.directory}/var/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:db]
command=couchdb -a ${parts.buildout.directory}/etc/couchdb.ini
stderr_logfile = ${parts.buildout.directory}/var/log/couchdb/couch-stderr.log
stdout_logfile = ${parts.buildout.directory}/var/log/couchdb/couch-stdout.log
environment=HOME=${parts.buildout.directory}/var


{% if options['develop'] %}
[program:webapp]
command = ${parts.buildout.directory}/bin/chaussette  --backend=gevent paste:${parts['openprocurement.api.ini'].output} --graceful-timeout=${options['graceful_timeout'] if 'graceful_timeout' in options else '60'}
numprocesses = ${options['numprocesses'] if 'numprocesses' in options else '1'}
environment=SANDBOX_MODE=True
stderr_logfile = ${parts.buildout.directory}/var/log/webapp.stdout.log
stdout_logfile = ${parts.buildout.directory}/var/log/webapp.stderr.log
{% end %}
